# rootfs/etc/services.d/mattermost/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Mattermost
# Take down the S6 supervision tree when Mattermost fails
# ==============================================================================

declare exit_code
readonly exit_code="${1}"

bashio::log.info "Mattermost service exited with code ${exit_code}"

if [[ "${exit_code}" -ne 0 ]] && [[ "${exit_code}" -ne 256 ]]; then
    bashio::log.warning "Halting add-on"
    /run/s6/basedir/bin/halt
fi

bashio::log.info "Service restart after closing"
