# rootfs/etc/services.d/mattermost/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Mattermost
# Runs the Mattermost server
# ==============================================================================

bashio::log.info "Starting Mattermost server..."

# Ensure proper ownership of mattermost directories
chown -R mattermost:mattermost /mattermost

# Change to mattermost directory
cd /mattermost || bashio::exit.nok "Could not change to Mattermost directory"

# Set environment variables
export MM_SERVICESETTINGS_SITEURL=""
export MM_LOGSETTINGS_ENABLECONSOLE=true
export MM_LOGSETTINGS_CONSOLELEVEL=INFO

# Execute Mattermost directly (running as root is fine for Home Assistant addons)
exec /mattermost/server/bin/mattermost